# deploy/kinoweekend.ini
[program:kinoweekend]
directory = /home/bezrukov/Envs/kinoweekend/source
user = bezrukov
command = /home/bezrukov/Envs/kinoweekend/source/gunicorn.sh
stdout_logfile = /home/bezrukov/Envs/kinoweekend/logs/supervisor_stdout.log
stderr_logfile = /home/bezrukov/Envs/kinoweekend/logs/supervisor_stderr.log
